# crontab
# Refresh 3x/dia (BRT): 06:00, 12:00, 15:00 = 09:00, 15:00, 18:00 UTC
0 9,15,18 * * * /app/refresh.sh
